# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-backup/run
#!/command/with-contenv bashio
# ==============================================================================
# Start MariaDB-backup service
# ==============================================================================
MARIADB_FIFODIR=/tmp/export_homeassistant_database_fifodir

s6-mkfifodir ${MARIADB_FIFODIR}

while :
do
    s6-ftrig-wait ${MARIADB_FIFODIR} "e"
    export-homeassistant-database
    s6-ftrig-notify ${MARIADB_FIFODIR} "E"
done
