# mariadb/rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv bashio
# ==============================================================================
# Start crond service
# ==============================================================================

# Start crond
bashio::log.info "Starting crond"
crond -f
