# mariadb/rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start crond service
# ==============================================================================

(crontab -l; echo "*/5     *       *       *       *       run-parts /etc/periodic/5min" ) | grep . | crontab -

exec crond -f -d 9
