# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-core/run
#!/command/with-contenv bashio
# ==============================================================================
# Start MariaDB service
# ==============================================================================
MARIADB_DATA=$(</tmp/run/MARIADB_DATA)

# Start mariadb
bashio::log.info "Starting MariaDB"
mkdir -p /run/mysqld
exec mysqld --datadir="${MARIADB_DATA}" --user=root < /dev/null
