# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-core/run
#!/command/with-contenv bashio
# ==============================================================================
# Start MariaDB service
# ==============================================================================
MARIADB_DATA=$(</etc/MARIADB_DATA)

# Start mariadb
bashio::log.info "Starting MariaDB..."
mysqld --datadir="${MARIADB_DATA}" --user=root < /dev/null &
MARIADB_PID=$!

# Wait until DB is running
while ! mysql -e "" 2> /dev/null; do
    sleep 1
done
bashio::log.info "MariaDB started"

# Notify s6 about successful service start
echo "" >&3

wait "${MARIADB_PID}"
