# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-core/run
#!/command/with-contenv bashio
# ==============================================================================
# Start MariaDB service
# ==============================================================================
MARIADB_DATA=$(</etc/MARIADB_DATA)

# Start mariadb
bashio::log.info "Starting MariaDB..."
rm -f /etc/MARIADB_SHUTDOWN
mysqld_safe --datadir="${MARIADB_DATA}" --user=root < /dev/null

bashio::log.info "MariaDB stopped"
if bashio::fs.file_exists "/etc/MARIADB_SHUTDOWN"; then
    # workaround to avoid restarting the service after intended shutdown
    sleep 10
fi
