# mariadb/rootfs/etc/services.d/mariadb/run
#!/command/with-contenv bashio
# ==============================================================================
# Start MariaDB service
# ==============================================================================
MARIADB_DATA=$(</etc/MARIADB_DATA)

# Register stop
function stop_mariadb() {
echo $(date -Iseconds) BEFORE export >> /share/log
    export-homeassistant-database
echo $(date -Iseconds) AFTER export >> /share/log
    bashio::services.delete "mysql"
    mysqladmin shutdown
    # Successful exit, avoid wait exit status to propagate
    exit 0
}
trap "stop_mariadb" SIGTERM SIGHUP

# Start mariadb
bashio::log.info "Starting MariaDB"
mysqld_safe --datadir="${MARIADB_DATA}" --user=root < /dev/null
